# supervisord.conf
[supervisord]
nodaemon=true

[program:couchdb]
command=/usr/bin/couchdb
user=couchdb
directory=/var/lib/couchdb
environment=HOME=/var/lib/couchdb
stopsignal=INT
